# roles/whalewall/templates/whalewall.service.j2
# {{ ansible_managed }}
# Whalewall systemd service
# Generated by Ansible whalewall role

[Unit]
Description=Whalewall - Automate management of firewall rules for Docker containers
After=docker.service
Wants=docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/whalewall
Restart=always
RestartSec=5
User=docker
Group=docker

[Install]
WantedBy=multi-user.target 